# deployment/init.d/simple-utcd
#!/bin/bash
#
# simple-utcd    Simple UTC Daemon init script
#
# chkconfig: 2345 20 80
# description: Simple UTC Daemon provides UTC time coordinate services
# processname: simple-utcd
# pidfile: /var/run/simple-utcd/simple-utcd.pid
# config: /etc/simple-utcd/simple-utcd.conf
#
# Author: SimpleDaemons
# License: Apache 2.0

### BEGIN INIT INFO
# Provides: simple-utcd
# Required-Start: $network $remote_fs $syslog
# Required-Stop: $network $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Simple UTC Daemon
# Description: Simple UTC Daemon provides UTC time coordinate services
### END INIT INFO

# Source function library
. /etc/rc.d/init.d/functions

# Source networking configuration
. /etc/sysconfig/network

# Check that networking is up
[ "${NETWORKING}" = "no" ] && exit 6

# UTC daemon configuration
UTCD_NAME="simple-utcd"
UTCD_BIN="/usr/sbin/simple-utcd"
UTCD_CONF="/etc/simple-utcd/simple-utcd.conf"
UTCD_PID="/var/run/simple-utcd/simple-utcd.pid"
UTCD_LOCK="/var/lock/subsys/simple-utcd"
UTCD_LOG="/var/log/simple-utcd/simple-utcd.log"
UTCD_USER="utc"
UTCD_GROUP="utc"
UTCD_ARGS="-c ${UTCD_CONF}"

# Check if binary exists
[ -x ${UTCD_BIN} ] || exit 5

# Check if config exists
[ -f ${UTCD_CONF} ] || exit 6

# Create necessary directories
create_dirs() {
    mkdir -p /var/run/simple-utcd
    mkdir -p /var/log/simple-utcd
    mkdir -p /var/lib/simple-utcd
    chown ${UTCD_USER}:${UTCD_GROUP} /var/run/simple-utcd
    chown ${UTCD_USER}:${UTCD_GROUP} /var/log/simple-utcd
    chown ${UTCD_USER}:${UTCD_GROUP} /var/lib/simple-utcd
    chmod 755 /var/run/simple-utcd
    chmod 755 /var/log/simple-utcd
    chmod 755 /var/lib/simple-utcd
}

# Start the daemon
start() {
    echo -n "Starting ${UTCD_NAME}: "

    # Create necessary directories
    create_dirs

    # Check if already running
    if [ -f ${UTCD_PID} ]; then
        if kill -0 $(cat ${UTCD_PID}) 2>/dev/null; then
            echo "already running"
            return 1
        else
            rm -f ${UTCD_PID}
        fi
    fi

    # Start the daemon
    daemon --user ${UTCD_USER} --pidfile ${UTCD_PID} ${UTCD_BIN} ${UTCD_ARGS}
    RETVAL=$?

    if [ $RETVAL -eq 0 ]; then
        echo "OK"
        touch ${UTCD_LOCK}
    else
        echo "FAILED"
    fi

    return $RETVAL
}

# Stop the daemon
stop() {
    echo -n "Stopping ${UTCD_NAME}: "

    if [ -f ${UTCD_PID} ]; then
        killproc -p ${UTCD_PID} ${UTCD_BIN}
        RETVAL=$?
        if [ $RETVAL -eq 0 ]; then
            rm -f ${UTCD_PID}
            rm -f ${UTCD_LOCK}
            echo "OK"
        else
            echo "FAILED"
        fi
    else
        echo "not running"
        RETVAL=0
    fi

    return $RETVAL
}

# Restart the daemon
restart() {
    stop
    sleep 2
    start
}

# Reload configuration
reload() {
    echo -n "Reloading ${UTCD_NAME}: "

    if [ -f ${UTCD_PID} ]; then
        if kill -HUP $(cat ${UTCD_PID}) 2>/dev/null; then
            echo "OK"
            return 0
        else
            echo "FAILED"
            return 1
        fi
    else
        echo "not running"
        return 1
    fi
}

# Check status
status() {
    if [ -f ${UTCD_PID} ]; then
        if kill -0 $(cat ${UTCD_PID}) 2>/dev/null; then
            echo "${UTCD_NAME} is running (pid: $(cat ${UTCD_PID}))"
            return 0
        else
            echo "${UTCD_NAME} is not running (stale pid file)"
            return 1
        fi
    else
        echo "${UTCD_NAME} is not running"
        return 3
    fi
}

# Check configuration
configtest() {
    echo -n "Testing ${UTCD_NAME} configuration: "

    if [ -f ${UTCD_CONF} ]; then
        if ${UTCD_BIN} -t -c ${UTCD_CONF} >/dev/null 2>&1; then
            echo "OK"
            return 0
        else
            echo "FAILED"
            return 1
        fi
    else
        echo "configuration file not found"
        return 1
    fi
}

# Show usage
usage() {
    echo "Usage: $0 {start|stop|restart|reload|status|configtest}"
    echo "  start       - Start the UTC daemon"
    echo "  stop        - Stop the UTC daemon"
    echo "  restart     - Restart the UTC daemon"
    echo "  reload      - Reload configuration"
    echo "  status      - Show daemon status"
    echo "  configtest  - Test configuration file"
    exit 1
}

# Main script logic
case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        restart
        ;;
    reload)
        reload
        ;;
    status)
        status
        ;;
    configtest)
        configtest
        ;;
    *)
        usage
        ;;
esac

exit $?
